# web/etc/crontab
23 4 * * * /usr/local/mushus-gcp-cfg/web/upgrade.sh